# supervisord.conf
[supervisord]
nodaemon=true

[program:activemq]
command=/opt/apache-activemq-5.14.3/bin/activemq start
user=root
